# moma_demos/piloting_demo/router/piloting-router.service
[Unit]
Description=Piloting router
After=network.target

[Service]
ExecStart=/home/smb/router.sh
RemainAfterExit=true
Type=oneshot

[Install]
WantedBy=multi-user.target
